# misc/quant.socket.in
# systemd unit file to allow quant to bind to 443/UDP

[Unit]
Description=Quant UDP/443

[Socket]
ListenDatagram=443

[Install]
RequiredBy=quant.service
