# misc/quant.service
[Unit]
Description=QUANT server

[Service]
ExecStart=/home/ubuntu/quant/misc/quant-server.sh
Restart=always

[Install]
WantedBy=default.target
